# data/artifacts/T1037_persistence.yaml
name: Custom.MITRE.T1037_persistence
type: CLIENT
description: persistence type of Logon Scripts with ID T1037
sources:
- query: "\nLET persistence_locations = SELECT * FROM Artifact.Mac.Persistence.LaunchAgents()\n\
    LIMIT 50\n\nSELECT \n    Name,\n    Label,\n    Program,\n    ProgramArguments,\n\
    \    RunAtLoad,\n    Source,\n    'T1037' as TechniqueID,\n    'Logon Scripts'\
    \ as TechniqueName\nFROM persistence_locations\n"
